# deploy/foodmood.service
[Unit]
Description=FoodMood — سامانه تغذیه سازمانی
Documentation=https://github.com/5emad/FoodMood
After=network-online.target mongod.service
Wants=network-online.target
Requires=mongod.service

[Service]
Type=simple
User=foodapp
Group=foodapp
WorkingDirectory=/opt/food
Environment=NODE_ENV=production
EnvironmentFile=/opt/food/.env
ExecStartPre=/opt/food/deploy/wait-for-mongo.sh
ExecStart=/usr/bin/node backend/server.js
Restart=always
RestartSec=5
StartLimitInterval=300
StartLimitBurst=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=foodmood

LogsDirectory=foodmood
LogsDirectoryMode=0750

# Hardening
NoNewPrivileges=true
PrivateTmp=false
ProtectSystem=strict
ProtectHome=true
Environment=HOME=/opt/food
Environment=FOOD_INSTALL_DIR=/opt/food
Environment=XDG_RUNTIME_DIR=/opt/food/.cache/pdf-runtime/run
ReadWritePaths=/opt/food /var/log/foodmood /tmp /run

[Install]
WantedBy=multi-user.target
